# deploy/linux/systemd/dgteam-backup.service
[Unit]
Description=DGTEAM Backup Live Release
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/srv/dgteam
EnvironmentFile=-/srv/dgteam/.env
LogsDirectory=dgteam
ExecStart=/srv/dgteam/deploy/linux/bin/dgteam-backup.sh
StandardOutput=append:/var/log/dgteam/backup.log
StandardError=append:/var/log/dgteam/backup.err.log
